# templates/thumbor.service.erb
[Unit]
Description=Thumbor Image Server <%= @inst %>
After=syslog.target

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/usr/local/bin/thumbor -c /etc/thumbor.conf -k /etc/thumbor.key -p <%= @inst %>
User=thumbor
Group=thumbor

StandardOutput=file:<%= @thumbor::path_log %>/thumbor.log
StandardError=file:<%= @thumbor::path_log %>r/thumbor.log
ReadWriteDirectories=-<%= @thumbor::path_log %>
SysLogIdentifier=thumbor

[Install]
WantedBy=multi-user.target


