# templates/thumbor.service.erb
[Unit]
Description=Thumbor Image Server <%= @inst %>
After=syslog.target

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/usr/local/bin/thumbor -c /etc/thumbor/thumbor-<%= @inst %>.conf -k /etc/thumbor/thumbor.key -p <%= @inst %>
User=thumbor
Group=thumbor

StandardOutput=file:<%= @plog %>/thumbor.log
StandardError=file:<%= @plog %>/thumbor.log
ReadWriteDirectories=-<%= @plog %>
SysLogIdentifier=thumbor

[Install]
WantedBy=multi-user.target


